# resources/docker/alpine/dnssetup
#!/sbin/openrc-run
# TODO: remove this file when we fix DNS

start() {
  echo "nameserver 8.8.8.8" > /etc/resolv.conf
  echo "127.0.0.1       localhost" >> /etc/hosts
}
